# deploy/slide-maker.service
[Unit]
Description=Slide Maker Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/slidemaker
ExecStart=/home/ubuntu/.bun/bin/bun run dist/server/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
